# infra/deploy/systemd/pv-backend.service
[Unit]
Description=PV Poverty Alleviation FastAPI backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/pv-poverty-alleviation/backend
EnvironmentFile=-/srv/pv-poverty-alleviation/backend/.env
Environment=PYTHONPATH=/srv/pv-poverty-alleviation/backend:/srv/pv-poverty-alleviation
Environment=TZ=Asia/Shanghai
ExecStart=/srv/pv-poverty-alleviation/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --proxy-headers
Restart=always
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
